# playbooks/files/supervisord.conf
[program:consul-template]
autostart=true
autorestart=true
command=consul-template -config /etc/consul-template/conf.d/